# CD_RaspberryPi2/etc/init.d/avr_daemon.sh~
#!/bin/sh
### BEGIN INIT INFO
# Provides: avr_daemon
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start:  2 3 5 S
# Default-Stop: 0 1 6
# Short-Description:  Enable RS-232C polling from AVR.
# Description: Scan RS-232C polling realtime from ttyUSB0 when AVR send unsigned char data.
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin
AVR_DAEMON_PATH=/usr/sbin/avr_daemon

case "$1" in
  start)
   	echo -n "Starting avr_daemon : "
   	daemon 
    	;;
  stop)
    	;;
  restart)
    	echo -n "Restart avr_daemon : "
    	$0 stop
    	$0 start
    	;;
    *)
  	echo "Usage : $0 {start | stop | restart}"
  	exit 1
esac
  	exit 0
